# redi301.conf
; /etc/supervisor/conf.d/redi301.conf

[program:redi301]
process_name=redi301
directory=/root/redi301
command=/root/redi301/redi301 -a 0.0.0.0:80
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/redi301.log
